# deploy/beleggeriq.service
# BeleggerIQ 2.0 — systemd unit
#
# Installatie:
#   sudo cp deploy/beleggeriq.service /etc/systemd/system/beleggeriq.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now beleggeriq
#
# Logs:      journalctl -u beleggeriq -f
# Restart:   sudo systemctl restart beleggeriq

[Unit]
Description=BeleggerIQ 2.0 (Next.js)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=beleggeriq
Group=beleggeriq
WorkingDirectory=/mnt/HC_Volume_105455257/apps/beleggeriq/current
EnvironmentFile=/mnt/HC_Volume_105455257/apps/beleggeriq/shared/.env.production
# Build-info (BIQ_GIT_SHA / BIQ_BUILD_TIME) — geschreven door deploy.sh.
# Optioneel ('-' prefix): de service start ook als 't bestand er niet is
# (eerste deploy, of bij rollback naar een release van vóór deze conventie).
EnvironmentFile=-/mnt/HC_Volume_105455257/apps/beleggeriq/shared/.env.build-info
ExecStart=/usr/bin/node /mnt/HC_Volume_105455257/apps/beleggeriq/current/node_modules/.bin/next start -p 3003
Restart=always
RestartSec=5
# Harden
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/mnt/HC_Volume_105455257/apps/beleggeriq /tmp
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
# Next.js logs gaan naar stdout/stderr → journald
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
